# systemd/ssh-telnet-eda.service
[Unit]
Description=SSH->Telnet proxy (ssh-telnet-eda)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/ssh-telnet-eda --config /etc/ssh-telnet-eda/config.toml
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
